# deploy/clipilot-registry.service
[Unit]
Description=CLIPilot Registry
After=network.target

[Service]
Type=simple
User=clipilot
Group=clipilot
WorkingDirectory=/opt/clipilot-registry
EnvironmentFile=/etc/clipilot-registry/env
ExecStart=/opt/clipilot-registry/registry
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
